# xmodmap.service




[Service]
/opt/bin/load_modkeymap.sh

[Install]
WantedBy=default.target



